# deploy/com.tj.quant-tracker.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly quant-tracker run (Sundays 18:00 PT): seed --refresh → screen → report.
  Install:
    cp deploy/com.tj.quant-tracker.weekly.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.tj.quant-tracker.weekly.plist
  Test once:  launchctl kickstart -k gui/$(id -u)/com.tj.quant-tracker.weekly
  Remove:     launchctl bootout gui/$(id -u)/com.tj.quant-tracker.weekly
              rm ~/Library/LaunchAgents/com.tj.quant-tracker.weekly.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tj.quant-tracker.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/user/dev/quant-tracker/bin/scheduled-run.sh</string>
        <string>weekly</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>18</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/user/dev/quant-tracker/logs/launchd.weekly.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/dev/quant-tracker/logs/launchd.weekly.err.log</string>
</dict>
</plist>
