# deploy/com.tj.quant-tracker.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  NOTE: replace every "/path/to/quant-tracker" below with this repo's absolute
  path on your machine (launchd requires absolute paths, no ~ expansion).
  Daily quant-tracker run (Mon-Fri 13:30 PT, after the 1pm PT / 4pm ET close):
  paper monitor → report. Install:
    cp deploy/com.tj.quant-tracker.daily.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.tj.quant-tracker.daily.plist
  Test once:  launchctl kickstart -k gui/$(id -u)/com.tj.quant-tracker.daily
  Remove:     launchctl bootout gui/$(id -u)/com.tj.quant-tracker.daily
              rm ~/Library/LaunchAgents/com.tj.quant-tracker.daily.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tj.quant-tracker.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/quant-tracker/bin/scheduled-run.sh</string>
        <string>daily</string>
    </array>
    <!-- An array of intervals = Mon(1) through Fri(5) at 13:30. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/path/to/quant-tracker/logs/launchd.daily.out.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/quant-tracker/logs/launchd.daily.err.log</string>
</dict>
</plist>
